# 10-viki.rules
SUBSYSTEM=="tty", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="5740", RUN+="/home/pos/ubuntu_init/check_viki_port.sh"


