# 12-arcus.rules
SUBSYSTEM=="tty", ATTRS{idProduct}=="903b", RUN+="/home/pos/ubuntu_init/check_arcus_port.sh"
